# dirt-scraper/bin/dirtd.service
[Unit]
Description=DIRT Daemon Service

[Service]
WorkingDirectory=/opt/dirtd
ExecStart=/usr/bin/java -Xms20m -Xmx200m \
        -cp "dirt-scraper.jar:lib/*" \
        -Ddirtd.numWorkers=8 \
        -Dconfig="cfg/db.ini" \
        -Dlog4j.configurationFile="cfg/log4j2.xml" \
        dev.pagefault.eve.dirtd.DirtDaemon

User=dirt
Type=simple
Restart=on-failure
RestartSec=10
SuccessExitStatus=143

[Install]
WantedBy=multi-user.target
